# data/friskby-controlpanel.service
[Unit]
Description=Friskby Controlpanel

[Service]
ExecStart=/usr/local/bin/fby_controlpanel
Restart=always

[Install]
WantedBy=multi-user.target
