# data/friskby.service
[Unit]
Description=Friskby logger
After=network-online.target
Requires=network-online.target

[Service]
ExecStart=/usr/local/bin/fby_manager
Restart=always
RestartSec=3600

[Install]
WantedBy=multi-user.target
